# Skill Management for EcoSynTech/scheduler/systemd/ecosyntech-scheduler.service
[Unit]
Description=EcoSynTech Scheduler Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.local/share/opencode/worktree/235aa3d499aaf6c96fd60e7c52a229e9ffd9b946/hidden-canyon
ExecStart=/usr/bin/node "/root/.local/share/opencode/worktree/235aa3d499aaf6c96fd60e7c52a229e9ffd9b946/hidden-canyon/Skill Management for EcoSynTech/scheduler/scheduler-runner.js"
Restart=on-failure
Environment=SCHEDULER_CONFIG=/root/.local/share/opencode/worktree/235aa3d499aaf6c96fd60e7c52a229e9ffd9b946/hidden-canyon/Skill Management for EcoSynTech/scheduler/config/scheduler.json
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
