# Skill-Management-for-EcoSynTech/scheduler/systemd/ecosyntech-scheduler.service
[Unit]
Description=EcoSynTech Scheduler Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ecosyntech-web
ExecStart=/usr/bin/node /root/ecosyntech-web/Skill-Management-for-EcoSynTech/scheduler/scheduler-runner.js
Restart=on-failure
RestartSec=10
Environment=SCHEDULER_CONFIG=/root/ecosyntech-web/Skill-Management-for-EcoSynTech/scheduler/config/scheduler.json
Environment=SCHEDULER_API_KEY=CHANGE_ME_TO_YOUR_API_KEY
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
